# nodejs-app-mongodb/app/backend/backend.service
[Unit]
Description=Durantic demo backend (PDF generator)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/app
EnvironmentFile=/etc/durantic/backend.env
ExecStart=/usr/bin/node /opt/app/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
